# package/base-files/files/etc/init.d/done
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=16
boot() {
	mount_root done
	rm -f /sysupgrade.tgz && sync

	#Reboot if usb subsystem doesnt start properly
	if [ -n "$(uci -q get system.usbcheck.path)" ]; then
		disable="$(uci -q get system.modem.disable)"
		path="$(uci -q get system.usbcheck.path)"
		if [ -z "$disable" ] || [ "$disable" -eq 0 2>/dev/null ] && \
			[ ! -d "$path" ]; then
			echo "USB system start up: FAILED... rebooting!"
			cnt=$(uci -q get system.usbcheck.counter)
			max=$(uci -q get system.usbcheck.max)
			if [ "$max" -gt "$cnt" ]; then
				cnt=$((cnt+1))
				uci -q set system.usbcheck.counter=$cnt
				uci commit system
				sync && reboot -f
			fi
		else
			echo "USB system start up: OK"
			uci -q set system.usbcheck.counter='0'
			uci commit system
		fi
	fi
}
