# package/base-files/files/etc/init.d/kmodloader
#!/bin/sh /etc/rc.common

START=39

start() {
	kmodloader /etc/modules-late.d/
}
