# package/base-files/files/etc/init.d/rclocal
#!/bin/sh /etc/rc.common

START=99
boot() {
	# process user commands
	[ -f /etc/rc.local ] && {
		sh /etc/rc.local
	}
}
