# package/base-files/files/etc/init.d/rclocal.openrc
#!/sbin/openrc-run

depend() {
	before done
}

start() {
	# process user commands
	[ -f /etc/rc.local ] && {
		sh /etc/rc.local
	}
}
