# package/base-files/files/etc/init.d/umount.openrc
#!/sbin/openrc-run

stop() {
	sync
	/bin/umount -a -d -r
}
