# package/teltonika/modbus_server/bin/files/modbus_server.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=98
STOP=99

BIN=modbus_server
CONFIG=modbus_server

rename_custom_regfile() {
	REGFILE_ON=$1
	REGFILE=$(uci -q get modbus_server.modbus.regfile)

	if [ "$REGFILE_ON" == "1" ]; then
		OLD_REGFILE="$(grep '# REGISTER_BLOCK_FILE' /etc/sysupgrade.conf 2>/dev/null | head -n1 | sed 's/ #.*//')"
		if [ "$OLD_REGFILE" != "$REGFILE" ]
		then
			rm -f "$OLD_REGFILE"
			sed -i '/# REGISTER_BLOCK_FILE/d' /etc/sysupgrade.conf 2>/dev/null
			echo "${REGFILE} # REGISTER_BLOCK_FILE" >> /etc/sysupgrade.conf
		fi
	fi
}

set_tcp_service() {
	local SECTION=$1
	ENABLED=$(uci -q get modbus_server.modbus.enabled)
	REGFILE_ON=$(uci -q get modbus_server.modbus.clientregs)

	rename_custom_regfile "$REGFILE_ON"

	if [ "$ENABLED" == "1" ]; then
		procd_open_instance "$SECTION"
		section_hash=$(uci show "$CONFIG"."$SECTION" | md5sum | head -c 32)
		procd_set_param data hash=$section_hash
		procd_set_param command $BIN -s "modbus" -H $section_hash
		procd_set_param respawn
		procd_close_instance
	fi
}

set_serial_service() {
	local SECTION=$1

	config_get ENABLED $SECTION enabled

	if [ -n "$ENABLED" ] && [ "$ENABLED" -eq 1 ]; then
		REGFILE_ON=$(uci -q get modbus_server.$SECTION.clientregs)
		rename_custom_regfile "$REGFILE_ON"

		procd_open_instance "$SECTION"
		section_hash=$(uci show "$CONFIG"."$SECTION" | md5sum | head -c 32)
		procd_set_param data hash=$section_hash
		procd_set_param command "$BIN" -s "$SECTION" -H $section_hash
		procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-60} ${respawn_retry:-0}
		procd_close_instance
	fi
}

start_service() {
	config_load "$CONFIG"

	set_tcp_service modbus

	config_foreach set_serial_service rtu_device
}

service_triggers() {
	procd_add_reload_trigger "modbus_server"
}
