# package/system/ubox/files/logd.openrc
#!/sbin/openrc-run

command="/sbin/logd"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
supervisor=supervise-daemon
respawn_delay=1
respawn_max=0

size="$(uci -q get system.system.log_size)"
[ "$size" -eq 0 ] && size=16
command_args="-S $size"

depend() {
	need ubusd system
}
