# package/feeds/vuci/vuci-app-mobile-utilities-api/files/etc/hotplug.d/gsm/1-new_modem
#!/bin/sh

([ "$EVENT_NAME" = "gsm.new_modem" ] || [ "$EVENT_NAME" = "gsm.modem_gone" ]) && /bin/trigger_vuci_routes_reload
