# package/system/urngd/files/urgnd.openrc
#!/sbin/openrc-run

command="/sbin/urngd"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	before boot
}
