# package/teltonika/azure_iothub/bin/files/azure_iothub.init
#!/bin/sh /etc/rc.common

START=99
STOP=99

PROG="/sbin/azure_iothub"
CONF="/etc/config/azure_iothub"
USE_PROCD=1

azure_instance() {
	config_get enabled "$1" "enabled"
	if [ "$enabled" = "1" ]; then
		procd_open_instance
		procd_set_param command "$PROG" -n "$1"
		procd_set_param file "$CONF"
		procd_close_instance
	fi
}

start_service() {
	config_load azure_iothub
	config_foreach azure_instance "azure_iothub"
}

service_triggers() {
	procd_add_reload_trigger "azure_iothub"
}
