# package/system/procd/files/procd.openrc
#!/sbin/openrc-run

command="/sbin/procd"
command_args="-h /etc/preinit.json"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
