# package/base-files/files/etc/hotplug.d/iface/91-trigger-vuci-routes-reload
case "$ACTION" in
"ifup" | "ifdown" | "ifupdate" | "free" | "create" | "add" | "remove")
	/bin/trigger_vuci_routes_reload -t "${ACTION} of ${INTERFACE} interface"
	;;
esac
